# data/key-mapper.rules
KERNEL=="sd*", RUN+="key-mapper-control --command autoload"
